# config/mysql/cron/mariadb.cron
0 3 * * * /etc/mysql/conf.d/cron/backup.sh
0 0 * * * /etc/mysql/conf.d/cron/log-rotate.sh
