# config/crowdsec/cron/crowdsec.cron
0 0 * * * /usr/local/cron/log-rotate.sh
