# config/mysql/mariadb.cron
0 5 1 */2 * /etc/mysql/conf.d/ssl.sh flush
0 3 * * * /etc/mysql/conf.d/backup.sh
0 0 * * * /etc/mysql/conf.d/log-rotate.sh